# CONTAINER/citadel_codeanalyst.container
[Unit]
Description=CITADEL repo app codeanalyst (citadel)

[Container]
ContainerName=citadel-codeanalyst
HostName=citadel-codeanalyst
Image=localhost/citadel-codeanalyst:latest
Network=host
AddCapability=NET_ADMIN
AddCapability=NET_RAW
AddDevice=/dev/net/tun
EnvironmentFile=/home/openclaw/saf/CODEANALYST/CONTAINER/citadel_codeanalyst.env
Volume=citadel-codeanalyst-var-lib-tailscale:/var/lib/tailscale:Z
# PublishPort lines are ignored in host mode

[Service]
Restart=always
RestartSec=5

[Install]
WantedBy=default.target
